# templates/default/etc/init/hook.conf.erb
start on runlevel [2345]
stop on runlevel [!2345]

setgid root
setuid root

exec /usr/local/<%= @name %> -port <%= @port %> -config <%= @conf_dir %> >> <%= @log %>
